# code/python_app.service
[Unit]

[Service]
ExecStart=/bin/sh -c "nohup python3 ../code/hello.py"