# code/python_app.service
[Unit]

[Service]
ExecStart=/bin/sh -c "nohup python3 /tmp/code/hello.py"
Restart=always
user=ec2-user